# pscmd.timer
[Unit]
Description=Таймер для pscmd

[Timer]
OnActiveSec=1

[Install]
WantedBy=default.target