# pscmd.timer
[Unit]
Description=Таймер для pscmd

[Timer]
OnActiveSec=30

[Install]
WantedBy=default.target