# pscmd.timer
[Unit]
Description=Таймер для pscmd

[Timer]
OnActiveSec=5

[Install]
WantedBy=default.target